# Scripts/launchd/studio.tmj.cardwise.cards-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
    <key>Label</key><string>studio.tmj.cardwise.cards-update</string>
    <key>ProgramArguments</key><array>
        <string>/bin/bash</string>
        <string>/Users/moltbot/CardWise-data/Scripts/update_cards.sh</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/moltbot/CardWise-data</string>
    <key>StartCalendarInterval</key><dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key><string>/Users/moltbot/CardWise-data/update.log</string>
    <key>StandardErrorPath</key><string>/Users/moltbot/CardWise-data/update.log</string>
    <key>EnvironmentVariables</key><dict>
        <key>PATH</key><string>/Users/moltbot/.local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict></plist>
